# systemd/user/redshift.service
[Unit]
Description=RED Shift screen

[Service]
Type=simple
ExecStart=/usr/bin/redshift

[Install]
WantedBy=xdm.target
